# factux/factux.service
[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/lampp/lampp start
ExecStop=/opt/lampp/lampp stop
